# etc/systemd/flowfuse.socket
[Unit]
Description=FlowFuse Socket
PartOf=flowfuse.service

[Socket]
ListenStream=3000

[Install]
WantedBy=sockets.target
